# crontab/pacron_np04dqm.cron
01 * * * * lxplus.cern.ch /afs/cern.ch/user/n/np04dqm/projects/dqmconfig/htcondor/pilotclean_np04dqm.sh
01 15 * * * lxplus.cern.ch /afs/cern.ch/user/n/np04dqm/projects/dqmconfig/htcondor/pilotpurge_np04dqm.sh
10 * * * * lxplus.cern.ch /afs/cern.ch/user/n/np04dqm/projects/dqmconfig/htcondor/pcalc_np04dqm.sh

