# crontab/pacron.cron
30 * * * * lxplus.cern.ch /afs/cern.ch/user/m/mxp/projects/dqmconfig/htcondor/pcalc.sh
55 * * * * lxplus.cern.ch /afs/cern.ch/user/m/mxp/projects/dqmconfig/services/olddata.sh
