# crontab/pacron_np04dqm.cron
1,11,21,31,41,51 * * * * lxplus.cern.ch /afs/cern.ch/user/n/np04dqm/projects/p3s/clients/tscan.sh -20 Proto
01 * * * * lxplus.cern.ch /afs/cern.ch/user/n/np04dqm/projects/dqmconfig/htcondor/pilotclean.sh
01 15 * * * lxplus.cern.ch /afs/cern.ch/user/n/np04dqm/projects/dqmconfig/htcondor/pilotpurge.sh
10 * * * * lxplus.cern.ch /afs/cern.ch/user/n/np04dqm/projects/dqmconfig/htcondor/pcalc_np04dqm.sh
01 12 * * * lxplus.cern.ch /afs/cern.ch/user/n/np04dqm/projects/dqmconfig/services/olddata.sh
