# crontab/pcron.cron
*/10 * * * * /afs/cern.ch/user/m/mxp/projects/dqmconfig/htcondor/pcalc.sh

