# crontab/pcron.cron
#!/bin/bash
/bin/date >> /tmp/crntst.out
